# files/services/matterbridge.service
[Unit]
Description=music bot for mumble

[Service]
ExecStart=bash -c 'cd /opt/programs/botamusique'

[Install]
WantedBy=default.target
